# debian/couchpotato.postinst
#!/bin/sh

mkdir -p /usr/lib/couchpotato

update-rc.d couchpotato defaults
